# files/root/service/unbound/run
#!/usr/bin/env sh

exec /usr/sbin/unbound -d -c /etc/unbound/unbound.conf 2>&1
